# screen-monitor-cleanup.service
[Unit]
Description=Screen Monitor - Screenshot Cleanup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/screen-monitor-cleanup
