# screen-monitor-archive-cleanup.service
[Unit]
Description=Screen Monitor Archive - Disk Rotation

[Service]
Type=oneshot
ExecStart=/usr/local/bin/screen-monitor-archive-cleanup
