# screen-monitor-sync.timer
[Unit]
Description=Screen Monitor - Sync Timer

[Timer]
OnCalendar=0/6:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
